# deploy/supervisord.ini
[program:totoro]
user=root
command=gunicorn -c deploy/prd_totoro_api.py run:app
autostart=true
autorestart=true
stopasgroup=true
directory=/var/www/totoro/server
environment=TOTORO_SECRET_KEY=""


[program:totoro_master]
user=root
command=gunicorn -c deploy/prd_totoro_admin.py run_admin:app
autostart=true
autorestart=true
stopasgroup=true
directory=/var/www/totoro/server
environment=TOTORO_SECRET_KEY=""
